# conf/wacom-inputattach@.service.in
[Unit]
Description=inputattach for Wacom ISDv4-compatible serial devices

[Service]
Type=simple
ExecStart=@BIN_PREFIX@/isdv4-serial-inputattach /dev/%I
Restart=on-success
